# automation/backup/com.user.gdrive-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.gdrive-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ivanvalenciaperez/Desktop/CLAUDE/BETA_PHASE/axkan-brain-v2/automation/backup/gdrive_backup.sh</string>
    </array>

    <!-- Run daily at 6:30 PM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <!-- Run missed jobs when system wakes up -->
    <key>StartInterval</key>
    <integer>86400</integer>

    <!-- Keep the job loaded even after it finishes -->
    <key>KeepAlive</key>
    <false/>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/Users/ivanvalenciaperez/gdrive_backup_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ivanvalenciaperez/gdrive_backup_stderr.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
